# roles/minio/templates/minio-upstart.conf
# /etc/init/minio.conf
 
description "minio"
author      "Bryan Liles <bl@do.co>"
 
start on (started caddy-minio)
stop on shutdown
 
# Automatically Respawn:
respawn
respawn limit 99 5
 
script
    USER="{{ minio_user }}"
    BIN="{{ minio_bin }}"    
    ADDR_OPTS="--address 127.0.0.1:9000 -C {{ minio_config_dir }}"
    OPTS="--min-free-disk 5% {{ minio_storage_dir }}"
    mkdir -p {{ minio_storage_dir }}
    mkdir -p {{ minio_config_dir }}
    chown $USER {{ minio_storage_dir }}
    chown $USER {{ minio_config_dir }}
    exec sudo -u $USER $BIN $ADDR_OPTS server $OPTS
end script